# aolab.rules
ATTRS{idVendor}=="0403", ATTRS{idProduct}=="6001", RUN+="/var/I1820.RPi/18.20-aolab.py"
